# deploy/wonderland-fetch.timer
[Unit]
Description=Weekly event scraper for wonderland

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
